# .config/systemd/user/power-profile@.service
[Unit]
Description=Set Power Profile to profile

[Service]
ExecStart=tuned-adm profile %i

[Install]
WantedBy=multi-user.target
