# .config/systemd/user/power-profile@desktop.timer
[Unit]
Description=Run desktop profile daily

[Timer]
OnCalendar=*-*-* 8:30:00

[Install]
WantedBy=timers.target
