# .config/systemd/user/power-profile@.service
[Unit]
Description=Set Power Profile to profile

[Service]
ExecStart=/usr/bin/tuned-adm profile %i

[Install]
WantedBy=default.target
