# .config/systemd/user/rclone@.service
[Unit]
Description=Rclone network mount for %i
Documentation=man:rclone(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
Environment="MOUNT_DIR=%h/mnt/%i"
Environment="REMOTE_PATH=/"
ExecStartPre=/usr/bin/mkdir -p ${MOUNT_DIR}
ExecStart=/home/linuxbrew/.linuxbrew/bin/rclone mount %i:${REMOTE_PATH} ${MOUNT_DIR} \
  --config=%h/.config/rclone/rclone.conf \
  --allow-other \
  --umask 0022 \
  --vfs-cache-mode full \
  --vfs-cache-max-size 50G \
  --vfs-cache-min-free-space 10G \
  --vfs-handle-caching 0s \
  --vfs-fast-fingerprint \
  --dir-cache-time 30s \
  --poll-interval 30s \
  --buffer-size 64M \
  --no-seek \
  --no-modtime \
  --fast-list \
  --async-read=true
ExecStop=/usr/bin/fusermount3 -uz ${MOUNT_DIR}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
