# .config/containers/systemd/power-profile@.service
[Unit]
Description=Set Power Profile

[Service]
ExecStart=tuned-adm profile %i

[Install]
WantedBy=multi-user.target
