# etc/systemd/system/power-profile@.service
[Unit]
Description=Set Power Profile

[Service]
Type=oneshot
ExecStart=tuned-adm profile %i

[Install]
WantedBy=multi-user.target
