# etc/systemd/system/power-profile@powersave.timer
[Unit]
Description=Switch to powersave profile at 00:30

[Timer]
OnCalendar=*-*-* 00:30:00

[Install]
WantedBy=timers.target
