# etc/systemd/system/var-mnt-nas.automount
[Unit]
Description=Automount for NFS NAS Mount

[Automount]
Where=/var/mnt/nas
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target
